# docker/basic-os/cron.d/tmp-clean
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Every hour, delete regular files under /tmp older than 3 days.
# Logs go to the container's stdout/stderr via PID 1 file descriptors.
0 * * * * root find /tmp -mindepth 1 -type f -mmin +4320 -delete >/proc/1/fd/1 2>/proc/1/fd/2
